# fruit_bot/fruit_bot.service
[Unit]
Description=Fruit Freshness Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/fruit_bot
EnvironmentFile=/home/ubuntu/fruit_bot/.env
ExecStart=/home/ubuntu/fruit_bot/venv/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
